# havoc_remote.service
[Unit]
Description=./HAVOC remote operator task

[Service]
Type=oneshot
ExecStart=/opt/havoc-remote/venv/bin/python3 /opt/havoc-remote/link.py

[Install]
WantedBy=multi-user.target